# astrbot@.service
[Unit]
Description=AstrBot Instance (%i)
After=network.target

[Service]
Type=simple
User=astrbot
Group=astrbot

# Runtime Environment
# UV_CACHE_DIR and UV_PROJECT_ENVIRONMENT are set by astrbotctl,
# but can be overridden here if needed.
Environment="ASTRBOT_ROOT=/var/lib/astrbot/%i"
Environment="ASTRBOT_SYSTEMD=1"
Environment="ASTRBOT_DESKTOP_CLIENT=0"
Environment="HOME=/var/lib/astrbot"
Environment="XDG_CACHE_HOME=/var/cache"
Environment="XDG_CONFIG_HOME=/etc"
Environment="XDG_DATA_HOME=/var/lib"
Environment="XDG_STATE_HOME=/var/lib"
Environment="UV_CACHE_DIR=/var/cache/astrbot/uv"
Environment="UV_PYTHON_INSTALL_DIR=/var/cache/astrbot/python"
Environment="UV_PROJECT_ENVIRONMENT=/var/cache/astrbot/venv-%i"

# Instance Configuration
EnvironmentFile=-/etc/astrbot/%i.conf

# Working directory
# We set this to the instance data directory for safety,
# though astrbotctl will cd to /opt/astrbot to find the project.
WorkingDirectory=/var/lib/astrbot

# Execution
# The wrapper script handles env setup and cd to /opt/astrbot
ExecStart=/usr/bin/astrbotctl run %i

Restart=always
RestartSec=10

# Security Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only

# Allow write access to necessary directories
ReadWritePaths=/var/lib/astrbot /var/cache/astrbot

[Install]
WantedBy=multi-user.target
